# ubuntu/lxc.conf
description "start services"

start on startup

script
initctl emit filesystem
end script
